# conf/supervisor.conf
[program:lecool]
command=/usr/bin/gunicorn application:app -b localhost:8009
directory=/home/ubuntu/sites/lecool/src/lecool/lecool
user=nobody
autostart=true
autorestart=true
